# pkg/jonnystoten/jonnystoten-drake.install
post_install() {
    systemctl enable --now lenovo_fix
}
